# webapp.service
[Unit]
Description=CSYE6225 Web Application
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=csye6225
Group=csye6225
WorkingDirectory=/opt/webapp
Environment="NODE_ENV=production"
EnvironmentFile=/opt/webapp/.env
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=webapp

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target